# src/terok_dbus/resources/systemd/terok-dbus.service
# SPDX-FileCopyrightText: 2026 Jiri Vyskocil
# SPDX-License-Identifier: Apache-2.0
#
# systemd user service for the terok-dbus Shield1 clearance hub.
#
# The {{BIN}} placeholder is resolved at install time by the install
# command (e.g. via sys.executable -m terok_dbus._cli) so the unit
# works regardless of PATH — matching the pattern used by
# terok-credential-proxy and terok-gate.
#
# Install:
#   terok-dbus install-service
#   — or manually —
#   sed "s|{{BIN}}|$(which terok-dbus)|" contrib/systemd/terok-dbus.service \
#       > ~/.config/systemd/user/terok-dbus.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now terok-dbus
#
# Logs:
#   journalctl --user -u terok-dbus -f

[Unit]
Description=terok Shield1 clearance hub (verdict dispatcher + desktop notifier)
Documentation=https://github.com/terok-ai/terok-dbus
After=dbus.socket
Requires=dbus.socket

[Service]
Type=simple
ExecStart={{BIN}} serve
Restart=on-failure
RestartSec=5
# No filesystem, namespace, or syscall hardening directives here, on
# purpose.  The hub exists to apply verdicts by shelling out to
# `terok-shield`, which in turn drives `podman unshare nsenter nft ...`.
# Several unit-file knobs that *seem* safe for a user service clash
# with that call path:
#
#   * `PrivateTmp`, `ProtectHome=read-only`, `ProtectSystem=full`,
#     `PrivateDevices` — mangle the mount / capability view enough
#     that rootless podman's re-exec-into-NS_ROOTLESS fails with
#     "cannot re-exec process to join the existing user namespace".
#
#   * `RestrictAddressFamilies`, `SystemCallFilter`, and similar
#     seccomp-based restrictions — the kernel refuses to install a
#     seccomp filter on an unprivileged process unless `no_new_privs`
#     is set (see `seccomp(2)`).  So even `NoNewPrivileges=no` can't
#     win: systemd has to enable NNP anyway to make the filter
#     load, and NNP in turn makes SELinux refuse the
#     `unconfined_t → container_runtime_t` transition that
#     `podman unshare` triggers.  Every verdict click then fails with
#     an `nnp_transition` AVC on SELinux-enforcing hosts (Fedora's
#     default).
#
# The hub runs as the operator's own uid with no privilege escalation
# over what their shell already has, and verdicts route through the
# auditable `terok-shield allow|deny` CLI — so the directives were
# hardening against a threat model the hub doesn't have.  Keep this
# section empty until there's a documented win for a specific
# directive that *doesn't* clash with rootless podman.

[Install]
WantedBy=default.target
