# src/terok_dbus/resources/systemd/terok-dbus.service
# SPDX-FileCopyrightText: 2026 Jiri Vyskocil
# SPDX-License-Identifier: Apache-2.0
#
# systemd user service for the terok-dbus Shield1 clearance hub.
#
# The {{BIN}} placeholder is resolved at install time by the install
# command (e.g. via sys.executable -m terok_dbus._cli) so the unit
# works regardless of PATH — matching the pattern used by
# terok-credential-proxy and terok-gate.
#
# Install:
#   terok-dbus install-service
#   — or manually —
#   sed "s|{{BIN}}|$(which terok-dbus)|" contrib/systemd/terok-dbus.service \
#       > ~/.config/systemd/user/terok-dbus.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now terok-dbus
#
# Logs:
#   journalctl --user -u terok-dbus -f

[Unit]
Description=terok Shield1 clearance hub (verdict dispatcher + desktop notifier)
Documentation=https://github.com/terok-ai/terok-dbus
After=dbus.socket
Requires=dbus.socket

[Service]
Type=simple
ExecStart={{BIN}} serve
Restart=on-failure
RestartSec=5
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
# ProtectHome=read-only locks down /run/user/$UID too; the event ingester
# needs to create and unlink its unix socket there, so whitelist it.  %t
# expands to $XDG_RUNTIME_DIR for user units.
ReadWritePaths=%t
PrivateDevices=yes
RestrictAddressFamilies=AF_UNIX

[Install]
WantedBy=default.target
