# src/terok_dbus/resources/systemd/terok-dbus.service
# SPDX-FileCopyrightText: 2026 Jiri Vyskocil
# SPDX-License-Identifier: Apache-2.0
#
# systemd user service for the terok-dbus Shield1 clearance hub.
#
# The {{BIN}} placeholder is resolved at install time by the install
# command (e.g. via sys.executable -m terok_dbus._cli) so the unit
# works regardless of PATH — matching the pattern used by
# terok-credential-proxy and terok-gate.
#
# Install:
#   terok-dbus install-service
#   — or manually —
#   sed "s|{{BIN}}|$(which terok-dbus)|" contrib/systemd/terok-dbus.service \
#       > ~/.config/systemd/user/terok-dbus.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now terok-dbus
#
# Logs:
#   journalctl --user -u terok-dbus -f

[Unit]
Description=terok Shield1 clearance hub (verdict dispatcher + desktop notifier)
Documentation=https://github.com/terok-ai/terok-dbus
After=dbus.socket
Requires=dbus.socket

[Service]
Type=simple
ExecStart={{BIN}} serve
Restart=on-failure
RestartSec=5
# No filesystem or namespace hardening here on purpose.  The hub exists
# to apply verdicts by shelling out to `terok-shield`, which in turn
# runs `podman unshare nsenter nft ...`.  Rootless podman re-execs
# itself to attach the pre-existing user namespace (pause process under
# $XDG_RUNTIME_DIR/libpod/), and that re-exec breaks under the typical
# "nice" unit-file knobs — PrivateTmp, ProtectHome=read-only,
# ProtectSystem=full, NoNewPrivileges, PrivateDevices all variously
# mangle the mount or capability view enough that podman fails with
# "cannot re-exec process to join the existing user namespace".
# The hub runs as the operator's own uid with no privilege escalation
# over what their shell already has, and every verdict it applies
# still routes through the auditable terok-shield CLI, so the extra
# directives were hardening against a threat model the hub doesn't
# have.  Keep this empty until there's a documented win for a specific
# directive that *doesn't* break podman.
RestrictAddressFamilies=AF_UNIX

[Install]
WantedBy=default.target
