# src/terok_clearance/resources/systemd/terok-clearance-verdict.service
# terok-clearance-verdict-version: {{UNIT_VERSION}}
# SPDX-FileCopyrightText: 2026 Jiri Vyskocil
# SPDX-License-Identifier: Apache-2.0
#
# systemd user service for the terok-clearance verdict helper.
#
# Stateless varlink server (``org.terok.ClearanceVerdict1``) that
# accepts one method — ``Apply(container, dest, action)`` — from
# the hub and exec's ``terok-shield allow|deny``.  Split from the
# hub so the hub can take tighter hardening; this unit carries only
# what the shield ``podman unshare nsenter nft`` re-exec tolerates.
#
# The {{BIN}} placeholder is resolved at install time.
#
# Logs:
#   journalctl --user -u terok-clearance-verdict -f

[Unit]
Description=terok clearance verdict helper (terok-shield exec)
Documentation=https://github.com/terok-ai/terok-clearance

[Service]
Type=simple
ExecStart={{BIN}} serve-verdict
Restart=on-failure
RestartSec=5

NoNewPrivileges=yes
AmbientCapabilities=
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
KeyringMode=private
UMask=0077
RestrictAddressFamilies=AF_UNIX AF_NETLINK AF_INET AF_INET6
IPAddressDeny=any

[Install]
WantedBy=default.target
