# src/terok_clearance/resources/systemd/terok-clearance-notifier.service
# SPDX-FileCopyrightText: 2026 Jiri Vyskocil
# SPDX-License-Identifier: Apache-2.0
#
# terok-clearance-notifier-version: {{UNIT_VERSION}}
#
# systemd user service for the terok-clearance-notifier — the process
# that bridges clearance-hub events to freedesktop desktop
# notifications.  Separate from the hub and verdict units so headless
# hosts can run the core clearance flow without pulling in a desktop
# stack.
#
# The {{BIN}} placeholder is resolved at install time by
# terok_clearance.runtime.installer.install_notifier_service so the
# unit works regardless of PATH.
#
# Install:
#   terok-sandbox setup            # or terok setup / terok-executor setup
#   — or manually —
#   sed "s|{{BIN}}|$(which python3) -m terok_clearance.notifier.app|" \
#       resources/systemd/terok-clearance-notifier.service \
#       > ~/.config/systemd/user/terok-clearance-notifier.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now terok-clearance-notifier
#
# Logs:
#   journalctl --user -u terok-clearance-notifier -f

[Unit]
Description=terok clearance desktop notifier
Documentation=https://github.com/terok-ai/terok-clearance
# Needs the hub running before we can subscribe; ``Wants=`` (not
# ``Requires=``) so a crashed hub restarts us with a clean retry rather
# than forcibly tearing us down.
Wants=terok-clearance-hub.service
After=terok-clearance-hub.service
# Freedesktop Notifications lives on the session bus.
After=dbus.socket
Requires=dbus.socket

[Service]
Type=simple
ExecStart={{BIN}}
Restart=on-failure
RestartSec=5

NoNewPrivileges=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
KeyringMode=private
UMask=0077
IPAddressDeny=any
RestrictAddressFamilies=AF_UNIX

[Install]
WantedBy=default.target
