# src/terok_clearance/resources/systemd/terok-clearance-verdict.service
# SPDX-FileCopyrightText: 2026 Jiri Vyskocil
# SPDX-License-Identifier: Apache-2.0
#
# terok-clearance-verdict-version: {{UNIT_VERSION}}
#
# systemd user service for the terok-clearance verdict helper.
#
# Stateless varlink server (``org.terok.ClearanceVerdict1``) that
# accepts one method — ``Apply(container, dest, action)`` — from the
# hub and exec's ``terok-shield allow|deny``.  The verdict exec
# transitively runs ``podman unshare nsenter nft``, which requires
# unhardened user+mount namespaces (see the long comment block in
# the earlier hub unit for the full story); that's why this helper
# deliberately ships without hardening directives.  The hardening
# lives on the hub side.
#
# Logs:
#   journalctl --user -u terok-clearance-verdict -f

[Unit]
Description=terok clearance verdict helper (terok-shield exec)
Documentation=https://github.com/terok-ai/terok-clearance

[Service]
Type=simple
ExecStart={{BIN}} serve-verdict
Restart=on-failure
RestartSec=5

# ── Intentionally NO hardening directives ────────────────────────────
#
# This unit execs ``terok-shield`` which in turn execs ``podman
# unshare nsenter nft``.  Rootless podman re-execs itself to attach
# the pause process's user+mount namespace, and two kernel-level
# rules break that re-exec under any real hardening:
#
# 1. Seccomp requires ``no_new_privs``; once NNP is on, SELinux
#    refuses the ``unconfined_t → container_runtime_t`` transition
#    and every verdict AVC-denies.
# 2. Mount-namespace isolation breaks setns-across-userns — even a
#    single ``InaccessiblePaths=`` entry causes the re-exec to fail.
#
# The hub (``terok-clearance-hub.service``) *does* carry the full
# hardening stack; isolating the unhardened shield exec to this
# helper is the whole point of the split.

[Install]
WantedBy=default.target
