# src/terok_clearance/resources/systemd/terok-clearance-notifier.service
# SPDX-FileCopyrightText: 2026 Jiri Vyskocil
# SPDX-License-Identifier: Apache-2.0
#
# terok-clearance-notifier-version: {{UNIT_VERSION}}
#
# systemd user service for the terok-clearance-notifier — the process
# that bridges clearance-hub events to freedesktop desktop
# notifications.  Separate from the hub and verdict units so headless
# hosts can run the core clearance flow without pulling in a desktop
# stack.
#
# The {{BIN}} placeholder is resolved at install time by
# terok_clearance.runtime.installer.install_notifier_service so the
# unit works regardless of PATH.
#
# Install:
#   terok-sandbox setup            # or terok setup / terok-executor setup
#   — or manually —
#   sed "s|{{BIN}}|$(which python3) -m terok_clearance.notifier.app|" \
#       resources/systemd/terok-clearance-notifier.service \
#       > ~/.config/systemd/user/terok-clearance-notifier.service
#   systemctl --user daemon-reload
#   systemctl --user enable --now terok-clearance-notifier
#
# Logs:
#   journalctl --user -u terok-clearance-notifier -f

[Unit]
Description=terok clearance desktop notifier
Documentation=https://github.com/terok-ai/terok-clearance
# Needs the hub running before we can subscribe; ``Wants=`` (not
# ``Requires=``) so a crashed hub restarts us with a clean retry rather
# than forcibly tearing us down.
Wants=terok-clearance-hub.service
After=terok-clearance-hub.service
# Freedesktop Notifications lives on the session bus.
After=dbus.socket
Requires=dbus.socket

[Service]
Type=simple
ExecStart={{BIN}}
Restart=on-failure
RestartSec=5

# Hardening — mirrors the hub's profile.  The notifier is the
# operator-UI-facing process that crosses adversarial container input
# (block-event payloads under attacker influence) into the gnome-shell
# session bus, so it deserves the same lockdown as the hub itself.
#
# ``ProtectHome=read-only`` is the right level for the threat model:
# the adversary's tool is a markup-injection payload aimed at
# gnome-shell — sanitised at the ``notifications/_sanitize.py``
# boundary — not a kernel CAP exploit, so write-protecting ``%h`` is
# what matters.  An earlier iteration used ``ProtectHome=tmpfs`` plus
# ``BindReadOnlyPaths=%h/.local/share/pipx/venvs/terok`` to bring the
# venv back; on Fedora-Atomic-style hosts (``/home`` symlinked to
# ``/var/home``) systemd's ``%h`` resolution and ``ProtectHome``'s
# coverage of ``/home`` vs ``/var/home`` don't cooperate cleanly and
# the notifier could end up running with no importable Python at all
# — silently zero desktop popups.  Read-only is simpler,
# threat-equivalent, and works regardless of pipx venv location.
NoNewPrivileges=yes
AmbientCapabilities=
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
SystemCallArchitectures=native
KeyringMode=private
UMask=0077

ProtectClock=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ProtectHostname=yes
ProtectProc=invisible
ProcSubset=pid

PrivateDevices=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only

RestrictAddressFamilies=AF_UNIX
SystemCallFilter=@system-service
SystemCallFilter=~@resources @privileged
MemoryDenyWriteExecute=yes
PrivateNetwork=yes
IPAddressDeny=any


[Install]
WantedBy=default.target
